# bin/crontab
*/10 * * * * bash /work/microbiome/sw/hpc_scripts/bin/check_disk >/dev/null 2>&1
